# 99-power-profile.rules
SUBSYSTEM=="power_supply", RUN+="/usr/local/bin/power-profile-switch.sh"