# power-profile-switch.service
[Unit]
Description=Set power profile based on AC status
After=multi-user.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/power-profile-switch.sh

[Install]
WantedBy=multi-user.target